# linux-whisperer.service
[Unit]
Description=Linux Whisperer - Voice Recording and Transcription Service
After=network.target

[Service]
Type=simple
User=%u
ExecStart=/usr/bin/python3 %h/code/linux-whisperer/main.py
Restart=on-failure
Environment=DISPLAY=:0
Environment=XAUTHORITY=%h/.Xauthority
Environment=DBUS_SESSION_BUS_ADDRESS=unix:path=/run/user/%U/bus

[Install]
WantedBy=default.target 